# rootfs/k3s-bootstrap
#!/sbin/openrc-run
# LatticeVE k3s bootstrap (OpenRC, oneshot). Reads per-node configuration from
# the Firecracker Metadata Service (MMDS, V1 @ 169.254.169.254) which the
# controller publishes pre/post-boot, prepares the node, and writes the resolved
# k3s command to /run/k3s/run. The separate, supervised `k3s` service then runs
# that wrapper (respawning k3s if it crashes) — see the k3s init script.
#
# MMDS keys:
#   role           "server" | "agent"
#   token          cluster join token
#   callback_token single-use token authenticating the kubeconfig POST (cp-0 only)
#   server         join URL (https://<ip>:6443) — for agents and joining servers
#   init           "1" on the first control-plane node (embedded-etcd cluster-init)
#   cni            CNI name ("flannel" default → host-gw; anything else → none)
#   callback       controller URL the first server POSTs its kubeconfig to
#   callback_ca    base64 controller serving cert, to verify the callback's TLS
#   node_name      k8s node name (hostname)
#   tls_san        space-separated extra TLS SANs (LB VIP, public IP) for the cert
#   root_pw_hash   crypt(3) root password hash (optional; default passwordless)
#   ssh_keys       base64 newline-joined SSH public keys (optional; enables sshd)
description="LatticeVE k3s bootstrap (Firecracker MMDS)"
depend() { need net; before k3s; }

mmds() { wget -qO- "http://169.254.169.254/$1" 2>/dev/null | tr -d '"'; }

# apply_node_access sets an optional root password and installs optional SSH
# authorized keys (starting dropbear, key-only, when keys are present). Both
# default to off: passwordless root, no SSH listener.
apply_node_access() {
    local pwhash=$(mmds root_pw_hash)
    if [ -n "$pwhash" ]; then
        echo "root:$pwhash" | chpasswd -e 2>/dev/null || true
    fi
    local sshk=$(mmds ssh_keys)
    if [ -n "$sshk" ]; then
        mkdir -p /root/.ssh
        if echo "$sshk" | base64 -d > /root/.ssh/authorized_keys 2>/dev/null; then
            chmod 700 /root/.ssh; chmod 600 /root/.ssh/authorized_keys
            # Key-only SSH: -s disables password logins, -g disables root password
            # logins; root may still log in by key.
            echo 'DROPBEAR_OPTS="-s -g"' > /etc/conf.d/dropbear
            rc-service dropbear start 2>/dev/null || true
        fi
    fi
}

start() {
    ebegin "LatticeVE k3s bootstrap"

    # The base rootfs image is built small; the controller truncates the raw disk
    # up to the provisioned size before boot, but truncate alone doesn't grow the
    # ext4 filesystem. Online-resize root so containerd has the full disk. Capture
    # the real error on failure rather than swallowing it — an undersized root
    # otherwise only surfaces much later as containerd disk pressure.
    _resize2fs=""
    for _p in /sbin/resize2fs /usr/sbin/resize2fs /usr/bin/resize2fs; do
        [ -x "$_p" ] && _resize2fs="$_p" && break
    done
    if [ -n "$_resize2fs" ]; then
        if ! out=$("$_resize2fs" /dev/vda 2>&1); then ewarn "resize2fs: $out"; fi
    else
        ewarn "resize2fs not found — root filesystem not grown (rebuild image to include e2fsprogs-extra)"
    fi

    ip route add 169.254.169.254/32 dev eth0 2>/dev/null || true

    local role="" i=0
    while [ $i -lt 90 ]; do
        role=$(mmds role); [ -n "$role" ] && break
        sleep 2; i=$((i+1))
    done
    [ -z "$role" ] && { eend 1 "MMDS role unavailable"; return 1; }

    local token=$(mmds token) cb_token=$(mmds callback_token) server=$(mmds server) cni=$(mmds cni) cb=$(mmds callback)
    local nn=$(mmds node_name) init=$(mmds init) tlssan=$(mmds tls_san)

    # Hostname → meaningful k8s node name (default is "localhost" otherwise).
    if [ -n "$nn" ]; then hostname "$nn"; echo "$nn" > /etc/hostname; fi

    # Guarantee a working resolver: VPC gateway (dnsmasq) + public fallback.
    # udhcpc does not reliably populate resolv.conf in this minimal rootfs.
    local gw=$(ip route | awk '/default/{print $3; exit}')
    { [ -n "$gw" ] && echo "nameserver $gw"; echo "nameserver 8.8.8.8"; } > /etc/resolv.conf

    mkdir -p /etc/rancher/k3s /var/log /run/k3s

    apply_node_access

    # flannel host-gw works over the (host-level) VXLAN overlay without an in-guest
    # VXLAN module; any non-flannel CNI is left to the user (flannel disabled).
    local fb="--flannel-backend=host-gw"
    [ -n "$cni" ] && [ "$cni" != "flannel" ] && fb="--flannel-backend=none --disable-network-policy"

    local nnarg=""; [ -n "$nn" ] && nnarg="--node-name $nn"
    local sanargs=""; for s in $tlssan; do sanargs="$sanargs --tls-san $s"; done

    # Resolve the k3s command into the wrapper the supervised `k3s` service runs.
    # Lives on tmpfs (/run) and carries the join token, so keep it root-only.
    if [ "$role" = "server" ]; then
        # First server initialises embedded etcd; additional servers join it.
        local mode=""
        if [ -n "$init" ]; then mode="--cluster-init"
        elif [ -n "$server" ]; then mode="--server $server"
        fi
        cat > /run/k3s/run <<EOF
#!/bin/sh
exec /usr/local/bin/k3s server $mode --token "$token" --write-kubeconfig-mode 644 $fb $nnarg $sanargs
EOF
        # Phone the kubeconfig home (controller patches the endpoint + stores it).
        # The kubeconfig carries cluster-admin creds, so verify the callback's TLS
        # against the controller's published serving cert (callback_ca); fall back
        # to --no-check-certificate only if no cert was published (older controller).
        if [ -n "$cb" ]; then
          local wget_tls="--no-check-certificate"
          local cca=$(mmds callback_ca)
          if [ -n "$cca" ]; then
            mkdir -p /usr/local/share/ca-certificates
            if echo "$cca" | base64 -d > /usr/local/share/ca-certificates/lattice-controller.crt 2>/dev/null \
               && update-ca-certificates >/dev/null 2>&1; then
              wget_tls=""
            fi
          fi
          ( for j in $(seq 1 360); do
              if [ -f /etc/rancher/k3s/k3s.yaml ]; then
                wget -q $wget_tls --header="X-Cluster-Token: $cb_token" \
                  --post-file=/etc/rancher/k3s/k3s.yaml -O- "$cb" >/dev/null 2>&1 && break
              fi
              sleep 5
            done ) &
        fi
    else
        cat > /run/k3s/run <<EOF
#!/bin/sh
exec /usr/local/bin/k3s agent --token "$token" --server "$server" $nnarg
EOF
    fi
    chmod 700 /run/k3s/run
    eend 0
}
