# rootfs/k3s-bootstrap
#!/sbin/openrc-run
# LatticeVE k3s bootstrap (OpenRC, oneshot). Reads per-node configuration from
# the Firecracker Metadata Service (MMDS, V2 @ 169.254.169.254) which the
# controller publishes pre/post-boot, prepares the node, and writes the resolved
# k3s command to /run/k3s/run. The separate, supervised `k3s` service then runs
# that wrapper (respawning k3s if it crashes) — see the k3s init script.
#
# MMDS keys:
#   role           "server" | "agent"
#   token          cluster join token
#   callback_token single-use token authenticating the kubeconfig POST (cp-0 only)
#   server         join URL (https://<ip>:6443) — for agents and joining servers
#   init           "1" on the first control-plane node (embedded-etcd cluster-init)
#   cni            CNI name ("flannel" default → host-gw; anything else → none)
#   callback       controller URL the first server POSTs its kubeconfig to
#   node_name      k8s node name (hostname)
#   tls_san        space-separated extra TLS SANs (LB VIP, public IP) for the cert
#   root_pw_hash   crypt(3) root password hash (optional; default passwordless)
#   ssh_keys       base64 newline-joined SSH public keys (optional; enables sshd)
description="LatticeVE k3s bootstrap (Firecracker MMDS)"
depend() { need net; before k3s; }

mmds_token=""
mmds_token_expires=0
refresh_mmds_token() {
    local now=$(date +%s 2>/dev/null || echo 0)
    if [ -n "$mmds_token" ] && [ "$now" -lt "$mmds_token_expires" ]; then
        return 0
    fi
    local token
    token=$(curl -fsS -X PUT \
        -H "X-metadata-token-ttl-seconds: 21600" \
        "http://169.254.169.254/latest/api/token" 2>/dev/null || true)
    if [ -n "$token" ]; then
        mmds_token="$token"
        mmds_token_expires=$((now + 21000))
        return 0
    fi
    return 1
}

mmds() {
    local key="$1"
    refresh_mmds_token || true
    if [ -n "$mmds_token" ]; then
        curl -fsS -H "X-metadata-token: $mmds_token" "http://169.254.169.254/$key" 2>/dev/null | tr -d '"'
    else
        # Backward-compatible fallback for older hosts still configured with V1.
        wget -qO- "http://169.254.169.254/$key" 2>/dev/null | tr -d '"'
    fi
}

shell_quote() {
    printf "'"
    printf '%s' "$1" | sed "s/'/'\\\\''/g"
    printf "'"
}

valid_node_name() {
    [ -z "$1" ] || printf '%s' "$1" | grep -Eq '^[A-Za-z0-9]([A-Za-z0-9.-]{0,251}[A-Za-z0-9])?$'
}

valid_metadata_value() {
    # MMDS values are controller-owned, but reject newlines/control chars so they
    # cannot corrupt the generated /run/k3s/run shell wrapper.
    [ -z "$1" ] || ! printf '%s' "$1" | grep -q '[[:cntrl:]]'
}

valid_url() {
    [ -n "$1" ] && valid_metadata_value "$1" && printf '%s' "$1" | grep -Eq '^https://[^[:space:]]+$'
}

valid_san() {
    [ -n "$1" ] && printf '%s' "$1" | grep -Eq '^[A-Za-z0-9_.:-]+$'
}

# apply_node_access sets an optional root password and installs optional SSH
# authorized keys (starting dropbear, key-only, when keys are present). Both
# default to off: passwordless root, no SSH listener.
apply_node_access() {
    local pwhash=$(mmds root_pw_hash)
    if [ -n "$pwhash" ]; then
        echo "root:$pwhash" | chpasswd -e 2>/dev/null || true
    fi
    local sshk=$(mmds ssh_keys)
    if [ -n "$sshk" ]; then
        mkdir -p /root/.ssh
        if echo "$sshk" | base64 -d > /root/.ssh/authorized_keys 2>/dev/null; then
            chmod 700 /root/.ssh; chmod 600 /root/.ssh/authorized_keys
            # Key-only SSH: -s disables password logins, -g disables root password
            # logins; root may still log in by key.
            echo 'DROPBEAR_OPTS="-s -g"' > /etc/conf.d/dropbear
            rc-service dropbear start 2>/dev/null || true
        fi
    fi
}

start() {
    ebegin "LatticeVE k3s bootstrap"

    # The base rootfs image is built small; the controller truncates the raw disk
    # up to the provisioned size before boot, but truncate alone doesn't grow the
    # ext4 filesystem. Online-resize root so containerd has the full disk. Capture
    # the real error on failure rather than swallowing it — an undersized root
    # otherwise only surfaces much later as containerd disk pressure.
    _resize2fs=""
    for _p in /sbin/resize2fs /usr/sbin/resize2fs /usr/bin/resize2fs; do
        [ -x "$_p" ] && _resize2fs="$_p" && break
    done
    if [ -n "$_resize2fs" ]; then
        if ! out=$("$_resize2fs" /dev/vda 2>&1); then ewarn "resize2fs: $out"; fi
    else
        ewarn "resize2fs not found — root filesystem not grown (rebuild image to include e2fsprogs-extra)"
    fi

    ip route add 169.254.169.254/32 dev eth0 2>/dev/null || true

    local role="" i=0
    while [ $i -lt 90 ]; do
        role=$(mmds role); [ -n "$role" ] && break
        sleep 2; i=$((i+1))
    done
    [ -z "$role" ] && { eend 1 "MMDS role unavailable"; return 1; }
    case "$role" in
        server|agent) ;;
        *) eend 1 "invalid MMDS role: $role"; return 1 ;;
    esac

    local token=$(mmds token) server=$(mmds server) cni=$(mmds cni)
    local nn=$(mmds node_name) init=$(mmds init) tlssan=$(mmds tls_san)
    if [ -z "$token" ] || ! valid_metadata_value "$token"; then eend 1 "invalid or missing MMDS token"; return 1; fi
    if ! valid_node_name "$nn"; then eend 1 "invalid MMDS node_name: $nn"; return 1; fi
    if [ -n "$cni" ] && ! printf '%s' "$cni" | grep -Eq '^[A-Za-z0-9_.-]+$'; then
        eend 1 "invalid MMDS cni: $cni"; return 1
    fi
    if [ "$role" = "agent" ] && ! valid_url "$server"; then
        eend 1 "invalid or missing MMDS server for agent"; return 1
    fi
    if [ "$role" = "server" ] && [ -z "$init" ]; then
        if ! valid_url "$server"; then eend 1 "invalid or missing MMDS server for joining control-plane"; return 1; fi
    fi

    # Hostname → meaningful k8s node name (default is "localhost" otherwise).
    if [ -n "$nn" ]; then hostname "$nn"; echo "$nn" > /etc/hostname; fi

    # Guarantee a working resolver: VPC gateway (dnsmasq) + public fallback.
    # udhcpc does not reliably populate resolv.conf in this minimal rootfs.
    local gw=$(ip route | awk '/default/{print $3; exit}')
    { [ -n "$gw" ] && echo "nameserver $gw"; echo "nameserver 8.8.8.8"; } > /etc/resolv.conf

    mkdir -p /etc/rancher/k3s /var/log /run/k3s

    apply_node_access

    # flannel host-gw works over the (host-level) VXLAN overlay without an in-guest
    # VXLAN module; any non-flannel CNI is left to the user (flannel disabled).
    local fb="--flannel-backend=host-gw"
    [ -n "$cni" ] && [ "$cni" != "flannel" ] && fb="--flannel-backend=none --disable-network-policy"

    local nnarg=""; [ -n "$nn" ] && nnarg="--node-name $(shell_quote "$nn")"
    local sanargs=""
    for s in $tlssan; do
        if ! valid_san "$s"; then eend 1 "invalid MMDS tls_san: $s"; return 1; fi
        sanargs="$sanargs --tls-san $(shell_quote "$s")"
    done
    local token_arg=$(shell_quote "$token")

    # Resolve the k3s command into the wrapper the supervised `k3s` service runs.
    # Lives on tmpfs (/run) and carries the join token, so keep it root-only.
    if [ "$role" = "server" ]; then
        # First server initialises embedded etcd; additional servers join it.
        local mode=""
        if [ -n "$init" ]; then mode="--cluster-init"
        elif [ -n "$server" ]; then mode="--server $(shell_quote "$server")"
        fi
        cat > /run/k3s/run <<EOF
#!/bin/sh
exec /usr/local/bin/k3s server $mode --token $token_arg --write-kubeconfig-mode 644 $fb $nnarg $sanargs --disable=traefik --disable=local-path --disable=metrics-server --disable=servicelb
EOF
    else
        local server_arg=$(shell_quote "$server")
        cat > /run/k3s/run <<EOF
#!/bin/sh
exec /usr/local/bin/k3s agent --token $token_arg --server $server_arg $nnarg
EOF
    fi
    chmod 700 /run/k3s/run
    eend 0
}
